# templates/supervisor.conf.erb
[program:<%= @app_name %>]
command=<%= @command %>
directory=<%= @directory %>
user=<%= @user %>
autostart=true
autorestart=true
redirect_stderr=True
<%- if @numprocs -%>
numprocs=<%= @numprocs %>
process_name=%(process_num)02d
<%- end -%>
<%- if @startsecs -%>
startsecs=<%= @startsecs %>
<%- end -%>
<%- if @stopwaitsecs -%>
stopwaitsecs=<%= @stopwaitsecs %>
<%- end -%>
<%- if @priority -%>
priority=<%= @priority %>
<%- end -%>
<%- if @environment %>
environment=<% @environment.each_pair do |key,value| -%><%=key%>="<%=value%>",<% end -%>
<%- end %>
<%- if @stopsignal %>
stopsignal=<%= @stopsignal %>
<%- end %>
<%- if @stopsignal %>
stopsignal=<%= @stopsignal %>
<%- end %>
<%- if @killasgroup%>
killasgroup=<%= @killasgroup%>
<%- end %>
<%- if @stopasgroup%>
stopasgroup=<%= @stopasgroup %>
<%- end %>
<% @additional_options.each do |key, value| %>
<%= key %>=<%= value %>
<% end %>


